# infrastructure/ansible/templates/laravel-scheduler.service.j2
[Unit]
Description=Animal Shelter Workshop - Laravel scheduler (single schedule:run pass)

[Service]
Type=oneshot
User=taufiq
WorkingDirectory=/home/taufiq/Animal-Shelter-Workshop
ExecStart=/usr/bin/php artisan schedule:run
